# src/guaucam-detector.service
[Unit]
Description=GuauCam noise detector (Telegram alerts) + web panel
After=guaucam-stream.service
StartLimitIntervalSec=0

[Service]
ExecStart=/usr/bin/python3 /opt/guaucam/src/noise_detector.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
